# slapd-setup.service
[Unit]

After=slapd.service

[Service]

ExecStart=/bin/bash /tmp/start.sh
RemainAfterExit=yes
Type=oneshot

[Install]

WantedBy=multi-user.target
